# deploy/rental-contract.service
[Unit]
Description=Lewis' Houses — Rental Property Management
After=network.target

[Service]
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/RentalContract
Environment="PATH=/home/ubuntu/RentalContract/venv/bin"
Environment="SECRET_KEY=change-this-to-a-random-secret-key"
ExecStart=/home/ubuntu/RentalContract/venv/bin/gunicorn -w 2 -b 0.0.0.0:5000 "app:create_app()"
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
